# agg.conf
[program:attempted_dot]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py attempteddot
autoreload=True
autorestart=True

[program:createshuffleorder]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py createshuffleorder
autoreload=True
autorestart=True

[program:copypooledquestions]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py copypooledquestions
autoreload=True
autorestart=True

[program:skillreport]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py skillreportupdate
autoreload=True
autorestart=True

[program:submission_details]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py submissiondetailsworker
autoreload=True
autorestart=True

[program:friends_image_fetch]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py invitefriendsgoogleimagefetch
autoreload=True
autorestart=True

[program:profile_completeness]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py saveprofilecompleteness
autoreload=True
autorestart=True

[program:transmissiondelete]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute satellite.consumer.delete_transmission
autoreload=True
autorestart=True

[program:transmissioninit]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute satellite.consumer.initialize_transmission
autoreload=True
autorestart=True

[program:signup_source]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py savesignupsource
autoreload=True
autorestart=True

[program: practice_section_elements_intro_mail]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py practicesectionelementsmail
autoreload=True
autorestart=True

[program: practice_section_elements_intro_mail]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py practicesectionelementsmail
autoreload=True
autorestart=True

[program:updateaggregatepenta]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py updateaggregatepenta
autoreload=True
autorestart=True

[program:update_analytics_worker]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py updateanalytics
autoreload=True
autorestart=True

[program:parse_resume]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py parseresume
autoreload=True
autorestart=True

[program:live_signup_map]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py live_user_signup
autoreload=True
autorestart=True

[program:user_labs]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute lab.worker.activate_user_labs
autoreload=True
autorestart=True

[program:avatar_error_handler]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py avatar_error_handler
autoreload=True
autorestart=True

[program:update_mparticipation_worker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py updatemparticipationworker
autoreload=True
autorestart=True

[program:push_init_submission]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py pushinitsubmission
autoreload=True
autorestart=True

[program:push_evaluated_submission]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py pushevaluatedsubmission
autoreload=True
autorestart=True

[program:update_thirdparty_account]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py update_thirdparty_account
autoreload=True
autorestart=True

[program:create_bulk_mcq]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py create_bulk_mcq
autoreload=True
autorestart=True

[program:updateshuffleonproblemsave]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py updateshuffleonproblemsave
autoreload=True
autorestart=True

[program:updatedisplayedproblemcount]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py updatedisplayedproblemcount
autoreload=True
autorestart=True

[program:anatomy_parent_updated]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py problemanatomyparentupdated
numprocs=5
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:fetch_source_from_repo]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py fetch_source_from_repo
autoreload=True
autorestart=True

[program:sprints_email_communicator]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py sprints_email_communicator
autoreload=True
autorestart=True

[program:update_user_geolocation]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py geolocationupdateconsumer
autoreload=True
autorestart=True

[program:create_update_invite_log]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py create_update_invite_log
autoreload=True
autorestart=True

[program:recommend_event_android_notifier]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py recommendeventworker
autoreload=True
autorestart=True

[program:not_found_error_worker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute hackathon.worker.consume_not_found_errors
autoreload=True
autorestart=True

[program:imageprocessconsumer]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py imageprocessconsumer
autoreload=True
autorestart=True

[program: capture_intercom_events]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py capture_intercom_events
autoreload=True
autorestart=True

[program: cleantestcase]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute problems.queue_worker.clean_testcase
autoreload=True
autorestart=True

[program: yearlyuserstats]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py yearlyuserstats
autoreload=True
autorestart=True

[program:registration_importer]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py registration_importer
autoreload=True
autorestart=True

[program:consume_sprint_events]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py consume_sprint_events
autoreload=True
autorestart=True

[program:update_project_score]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py updateprojectscore
autoreload=True
autorestart=True

[program:pusher]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py pusherworker
numprocs=15
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:match_scheduler_worker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py match_scheduler_worker
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:match_hoster_cron]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py match_hoster_cron
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:game_scheduler_worker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py game_scheduler_worker
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:rating_scheduler_cron]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py rating_scheduler_cron
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:bot_rating_worker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py bot_rating_worker
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:automated_emails]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py automated_emails
numprocs=5
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:onboarding_emails]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py onboarding_emails
numprocs=4
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:bulk_emails]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py send_bulk_emails
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:moderator_send_mails]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py marketing_emails
autoreload=True
autorestart=True

[program:blackops_emails]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute blackops.mails.queue.consume_message
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:challenge_anomaly_pipeline]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute anomaly.challenges.consumer.consume_event_ended
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:user_code_cache]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py invusercodecacheworker
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:uploadlogsandstopemulator]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py uploadlogsandstopemulator
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:update_topic_analytics_on_submission]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py updateanatomyanalyticsonsubmission
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:updatecommenttreeupvote]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py updatecommenttreeupvote
autoreload=True
autorestart=True

[program:rejudge_submission_worker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py rejudgesubmissionworker
autoreload=True
autorestart=True

[program: first_problem_solved]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py firstproblemsolved
autoreload=True
autorestart=True

[program: first_problem_solved]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py firstproblemsolved
autoreload=True
autorestart=True

[program:participation_ended]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py participationendeddetailsworker
autoreload=True
autorestart=True

[program:skillreport]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py skillreportupdate
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:update_analytics_worker]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py updateanalytics
autoreload=True

[program:rpupdate]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py analyticsupdate participation
numprocs=10
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:reupdate]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py analyticsupdate event
numprocs=6
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:pupdate]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py updatepracticereportasync
process_name=%(program_name)s_%(process_num)02d
numprocs=6
autoreload=True
autorestart=True

[program:plagiarism_report]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py plagiarismreportupdate
autoreload=True
autorestart=True

[program:email_delivery_lowest]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py email_delivery_worker lowest
process_name=%(program_name)s_%(process_num)02d
numprocs=8
autoreload=True
autorestart=True

[program:email_delivery_low]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py email_delivery_worker low
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autoreload=True
autorestart=True

[program:email_delivery_medium]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py email_delivery_worker medium
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autoreload=True
autorestart=True

[program:email_delivery_high]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py email_delivery_worker high
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autoreload=True
autorestart=True

[program:email_delivery_highest]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py email_delivery_worker highest
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autoreload=True
autorestart=True

[program:post_signal_worker]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py postsignalworker
autoreload=True
autorestart=True

[program:crud_calendar_events]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py crudcalendarevents
autoreload=True
autorestart=True

[program:oxygenreportworker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py oxygenreportworker
autoreload=True
autorestart=True

[program:hadoop_callback]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute problems.queue_worker.consume_mr_job_evaluated_queue
autoreload=True
autorestart=True

[program:pig_syantx_checker_callback]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute problems.worker.pig_syntax_checker_worker
autoreload=True
autorestart=True

[program: fight_club_submissions]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute fight_club.worker.process_evaluated_fight_submissions
autoreload=True
autorestart=True

[program: fight_club_bot_fights]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute fight_club.worker.process_bot_fights
autoreload=True
autorestart=True

[program: fight_club_bot_submissions]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute fight_club.worker.bot_submission_worker
autoreload=True
autorestart=True

[program: fight_club_mark_finished_fights]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute fight_club.worker.mark_finished_fights
autoreload=True
autorestart=True

[program: fight_club_update_user_levels]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute fight_club.worker.fight_club_update_levels
autoreload=True
autorestart=True

[program:hereminder_auto_subscribe]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute he_reminders.consume.consume_reminder_subscription_queue
autoreload=True
autorestart=True

[program:reminder_worker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute he_reminders.consume.consume_reminder_queue
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:smsworker]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute he_sms.consume.consume_sms_queue
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:update_circles_report]
command = python /home/careerstack/webapps/django/mycareerstack/manage.py execute circles.analytics_worker.update_circles_report
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autoreload=True
autorestart=True

[program:watcher]
command=python /home/careerstack/webapps/django/mycareerstack/manage.py execute watcher.worker.consume_watcher_queue
autoreload=True
autorestart=True
